# docker/graphite/etc/supervisor.d/statsd.ini
[program:statsd]
directory=/opt/statsd
command=/usr/bin/node stats.js config.js
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
